# deploy/supervisor/develop/rpress-web.conf
[program:rpress-web]
directory = %(ENV_RPRESS_HOME)s/rpress
command = %(ENV_RPRESS_HOME)s/venv/bin/gunicorn -b 127.0.0.1:5000 rpress:application

redirect_stderr = true
stdout_logfile = %(ENV_RPRESS_HOME)s/logs/rpress-web.log
